# systemd/user/change_wallpaper.timer
[Unit]
Description=Timer to change wallpapers every hour

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
